# docker.supervisor.conf
[program:docker]
command=/u01/wrapdocker
