# jenkins.supervisor.conf
[program:jenkins]
command=java -jar jenkins.war  --prefix=/jenkins
user=www-data
environment=HOME="/u01/jenkins_home",USER="www-data",JENKINS_HOME="/u01/jenkins_home"
